# Nodejs/service/supervisor.conf
<<<<<<< HEAD
[program:site]
command=node .
directory=/home/njs/site
user=njs
=======
[program:site]
command=node .
directory=/home/njs/site
user=njs
>>>>>>> 199b9bb2a68ef8300ea783cd497951e2c7cf8fc9
